# com.nedzib.pr_watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.nedzib.pr_watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/nedzibsastoque/Documents/Personal/review_inator/pr_watcher.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/nedzibsastoque/Documents/Personal/review_inator</string>

  <!-- Reiniciar automáticamente si el proceso muere -->
  <key>KeepAlive</key>
  <true/>

  <!-- Arrancar al hacer login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Logs -->
  <key>StandardOutPath</key>
  <string>/Users/nedzibsastoque/Documents/Personal/review_inator/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/nedzibsastoque/Documents/Personal/review_inator/launchd.err.log</string>

  <!-- Variables de entorno necesarias para gh, git, tmux, claude -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/nedzibsastoque</string>
  </dict>

</dict>
</plist>
